# .service
[Unit]
Description={{ app_name }}.service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User={{ ansible_ssh_user }}
ExecStart=/usr/bin/sudo /usr/bin/env python /home/{{ ansible_ssh_user }}/{{ app_name }}/app.py

[Install]
WantedBy=multi-user.target